# build/deb/pack-deb/DEBIAN/prerm
#!/bin/sh

set -e

svc=iothub

systemctl stop ${svc}.service
systemctl disable ${svc}.service

exit 0
